# debian-live/config/includes.chroot/etc/skel/.config/autostart/aulalliure-first-run.desktop
[Desktop Entry]
Type=Application
Name=AulaLliure Browser
Exec=firefox-esr --new-window --kiosk https://bibliojocs.edutictac.es/
Terminal=false
X-Moksha-Autostart=true
X-GNOME-Autostart-enabled=true
