# debian-live/config/includes.binary/aulalliure/etc/xdg/autostart/aulalliure-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=AulaLliure Wallpaper
Exec=/bin/true
X-GNOME-Autostart-enabled=false
Hidden=true
NoDisplay=true
NoDisplay=true
